# x010-supervisord/supervisor.conf
[program:go-name]
command=/usr/local/go/bin/go run main.go
autostart=true
autorestart=true
startretries=3

